# root/etc/s6-overlay/s6-rc.d/svc-yacht/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 3000" \
        cd /app/ s6-setuidgid abc /usr/bin/node server/index.mjs